# crontab
30 5 * * * /kinslayer/ImportPlayerLogs
# remember to end this file with an empty new line
